# example/11-evolutionary-meeting/crontab
# Evolutionary Meeting Bot - CEO単独会議
# 毎時00分に実装 + 振り返りを1セットで実行
#
# 分 時 日 月 曜日 コマンド

0 * * * * /app/scripts/meeting.sh >> /logs/cron.log 2>&1
